# supervisor/tests/cutomer_logger/supervisord.conf
[unix_http_server]
file=/dev/shm/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/dev/null
pidfile=/dev/shm/supervisord.pid
stderr_logfile=/dev/null
stdout_logfile=/dev/null
childlogdir=/opt/supervisor/logs/
logfile_maxbytes=0

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock ; use a unix:// URL  for a unix socket

[program:messages]
command=/opt/supervisor/supervisor/tests/cutomer_logger/messages
stderr_logfile=/opt/supervisor/logs/messages-stderr.log
stdout_logfile=/opt/supervisor/logs/messages-stdout.log
when=H

[program:messages1]
command=/opt/supervisor/supervisor/tests/cutomer_logger/messages1
stderr_logfile=/opt/supervisor/logs/messages1-stderr.log
stdout_logfile=/opt/supervisor/logs/messages1-stdout.log
when=D

[program:messages2]
command=/opt/supervisor/supervisor/tests/cutomer_logger/messages2
stderr_logfile=/opt/supervisor/logs/messages2-stderr.log
stdout_logfile=/opt/supervisor/logs/messages2-stdout.log
when=M

